# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:api]
command=bun apps/api/src/index.ts
directory=/app
user=chatty
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/tmp/api.stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/api.stderr.log
stderr_logfile_maxbytes=0

[program:web]
command=bun run start -p 3000
directory=/app/apps/web
user=chatty
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/tmp/web.stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/web.stderr.log
stderr_logfile_maxbytes=0
